# patches/etc/init/init-diaspora.conf
# init diaspora configuration

description "initializes diaspora on first start"
author "Scaleway <opensource@scaleway.com>"

start on (startup)

task

script
	/usr/local/sbin/oc-init-diaspora
end script
